# services/archiver.service
[Unit]
Description=Archier

[Service]
Type=oneshot
ExecStart=/usr/bin/archiver %i
